# systemd/rpicam-still-to-slack.service
[Unit]
Description=Capture still image and upload it to Slack
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pi/rpicam-still-to-slack
ExecStart=/home/pi/rpicam-still-to-slack/.venv/bin/python /home/pi/rpicam-still-to-slack/capture_to_slack.py
User=pi
Group=pi
